# files/plus/etc/uci-defaults/95-zram-swap
#!/bin/sh

# Peak-memory safety net for 1GB NanoPi R2S (plus/bypass).
# zram only compresses when pages are swapped; idle cost is near zero.

[ -x /etc/init.d/zram ] || exit 0
[ -f /etc/config/system ] || exit 0

# 256 MiB virtual swap is enough for short spikes without encouraging thrash.
uci set system.@system[0].zram_size_mb='256'

# ImmortalWrt prebuilt kmod-zram typically only ships LZO backends.
# LuCI (luci-app-zram-swap / system page) only offers: lzo | lz4 | zstd.
# Use "lzo" so the UI shows a selected value; "lzo-rle" works at runtime
# but is missing from the dropdown and appears as "-- 请选择 --".
algo='lzo'
if [ -r /sys/block/zram0/comp_algorithm ]; then
	# Prefer plain lzo when listed (matches LuCI). Fall back only if kernel
	# exposes neither lzo nor lzo-rle (unlikely).
	if grep -qw 'lzo' /sys/block/zram0/comp_algorithm 2>/dev/null \
		|| grep -q '\[lzo\]' /sys/block/zram0/comp_algorithm 2>/dev/null \
		|| grep -q 'lzo' /sys/block/zram0/comp_algorithm 2>/dev/null; then
		algo='lzo'
	elif grep -q 'lzo-rle' /sys/block/zram0/comp_algorithm 2>/dev/null; then
		algo='lzo-rle'
	fi
fi
uci set system.@system[0].zram_comp_algo="$algo"

uci set system.@system[0].zram_priority='100'
uci commit system

/etc/init.d/zram enable 2>/dev/null || true

exit 0
