# files/plus/etc/uci-defaults/95-zram-swap
#!/bin/sh

# Peak-memory safety net for 1GB NanoPi R2S (plus/bypass).
# zram only compresses when pages are swapped; idle cost is near zero.

[ -x /etc/init.d/zram ] || exit 0
[ -f /etc/config/system ] || exit 0

# 256 MiB virtual swap is enough for short spikes without encouraging thrash.
uci set system.@system[0].zram_size_mb='256'

# Current ImmortalWrt prebuilt kmod-zram only ships LZO backends.
# Prefer distro default lzo-rle; fall back to lzo if needed.
algo='lzo-rle'
if [ -r /sys/block/zram0/comp_algorithm ]; then
	if grep -q 'lzo-rle' /sys/block/zram0/comp_algorithm; then
		algo='lzo-rle'
	elif grep -q 'lzo' /sys/block/zram0/comp_algorithm; then
		algo='lzo'
	fi
fi
uci set system.@system[0].zram_comp_algo="$algo"

uci set system.@system[0].zram_priority='100'
uci commit system

/etc/init.d/zram enable 2>/dev/null || true

exit 0
