# files/common/etc/uci-defaults/98-r2s-network-device-migration
#!/bin/sh

changed=0

for section in $(uci show network 2>/dev/null | sed -n 's/^\(network\.[^.]*\)\.ifname=.*/\1/p'); do
	ifname="$(uci -q get "${section}.ifname")"
	[ -n "$ifname" ] || continue

	if [ -z "$(uci -q get "${section}.device")" ]; then
		uci set "${section}.device=$ifname"
	fi
	uci -q delete "${section}.ifname"
	changed=1
done

[ "$changed" = '1' ] && uci commit network

exit 0
