# files/bypass/etc/uci-defaults/96-bypass-sysupgrade
#!/bin/sh

# Migrate legacy whole-directory keep rules that break sysupgrade on R2S:
# boot partition is only 16MB; /etc/soho/ (geo+logs) must not be archived.

SYSUPGRADE_FILE="/etc/sysupgrade.conf"

touch "$SYSUPGRADE_FILE"

# Drop oversized / legacy directory keep rules (hooks keep essential files).
sed -i \
	-e '\|^/etc/soho/|d' \
	-e '\|^/etc/soho$|d' \
	"$SYSUPGRADE_FILE" 2>/dev/null || true

exit 0
